# bin/crontab
@midnight $HOME/bin/backup.sh >> /tmp/backup-${LOGNAME}.log 2>&1
